# lgswitch.service
[Unit]
Description=LG TV Auto-Switcher - USB Keyboard Monitor
Documentation=file:///home/robert/funspace/lgswitch/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/robert/funspace/lgswitch/venv/bin/python /home/robert/funspace/lgswitch/lgswitch.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Prevent rate limiting on restarts
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=default.target
