# lgswitch.service
[Unit]
Description=LG TV Auto-Switcher - USB Keyboard Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=INSTALL_PATH/venv/bin/python INSTALL_PATH/lgswitch.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Prevent rate limiting on restarts
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=default.target
